# config/supervisord.conf
[supervisord]
nodaemon=true
childlogdir=/var/log/lb_dashboard


[program:unicorn]
command=bundle exec unicorn_rails -E production -c ./config/unicorn.rb
